# root/etc/services.d/tail/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc tail -F "${CONFIG}/data/logs/CouchPotato.log"
